# dashboard/dashboard.conf
# put this file in ubuntu 12.04 /etc/init/dashboard.conf
# then to enable autocomplete sudo ln -s /etc/init/dashboard.conf /etc/init.d/dashboard
# sudo service dashboard start | stop
############################################################################################
start on (runlevel [2345])
exec /home/vagrant/dashboard/run.sh